# overlay/etc/cron.d/drush-mail-securityupdates
0   0 * * *    www-data            [ -x /usr/local/sbin/drush-mail-securityupdates ] && /usr/local/sbin/drush-mail-securityupdates
